# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=50MB
logfile_backups=3
loglevel=info
pidfile=/var/run/supervisord.pid

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

# ── PostgreSQL ──
[program:postgres]
command=/scripts/init-postgres.sh
user=root
autostart=true
autorestart=false
priority=10
stdout_logfile=/var/log/supervisor/postgres.log
stderr_logfile=/var/log/supervisor/postgres-err.log
startsecs=5

# ── Redis ──
[program:redis]
command=redis-server --save "" --appendonly no --bind 127.0.0.1 --port 6379
user=root
autostart=true
autorestart=true
priority=20
stdout_logfile=/var/log/supervisor/redis.log
stderr_logfile=/var/log/supervisor/redis-err.log
startsecs=3

# ── SearXNG ──
[program:searxng]
command=/scripts/start-searxng.sh
user=root
autostart=true
autorestart=true
priority=30
stdout_logfile=/var/log/supervisor/searxng.log
stderr_logfile=/var/log/supervisor/searxng-err.log
startsecs=5

# ── FastAPI ──
[program:fastapi]
command=/scripts/start-fastapi.sh
directory=/app
user=root
autostart=true
autorestart=true
priority=40
stdout_logfile=/var/log/supervisor/fastapi.log
stderr_logfile=/var/log/supervisor/fastapi-err.log
startsecs=10
environment=PYTHONPATH="/app"

# ── Nginx ──
[program:nginx]
command=nginx -g "daemon off;"
user=root
autostart=true
autorestart=true
priority=50
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx-err.log
startsecs=3
